# data/export/custom_systemd/master.target.erb
[Unit]
Wants=<%= process_names.join(' ') %>

[Install]
WantedBy=multi-user.target
